# ipset-nuclear/ipset-nuclear.timer
[Unit]
Description=Run ipset-nuclear periodically

[Timer]
OnUnitActiveSec=12h 1m
Persistent=true

[Install]
WantedBy=timers.target
